# sysroot/usr/lib/systemd/system/systemd-localed.service.d/locale-gen.conf
[Service]
# systemd-localed may run locale-gen which writes to /etc as well as to
# /usr/lib/locale. This change softens the service hardening a bit so
# both of these paths are writable.
ReadWritePaths=/usr/lib/locale/
