# sysroot/usr/lib/systemd/system/glamor-test.service
[Unit]
Description=Check for glamor

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/glamor_test.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
