# sysroot/usr/lib/udev/rules.d/99-libsane1.rules

ENV{DEVNAME}!="", ENV{libsane_matched}=="yes", RUN+="/bin/setfacl -m g:scanner:rw $env{DEVNAME}"
