# sysroot/usr/lib/systemd/user/gvfs-daemon.service
[Unit]
Description=Virtual filesystem service

[Service]
ExecStart=/usr/libexec/gvfsd
Type=dbus
BusName=org.gtk.vfs.Daemon
