# sysroot/usr/lib/systemd/user/org.gnome.SettingsDaemon.Datetime.service
[Unit]
Description=GNOME date & time service
CollectMode=inactive-or-failed
RefuseManualStart=true
RefuseManualStop=true

After=gnome-session-initialized.target

# Requisite/PartOf means the dependency is not loaded automatically.
# The ordering here also implies Before=gnome-session.target
Requisite=org.gnome.SettingsDaemon.Datetime.target
PartOf=org.gnome.SettingsDaemon.Datetime.target
Before=org.gnome.SettingsDaemon.Datetime.target



[Service]
Slice=session.slice
Type=dbus
ExecStart=/usr/libexec/gsd-datetime
Restart=on-failure
BusName=org.gnome.SettingsDaemon.Datetime
TimeoutStopSec=5
# We cannot use OnFailure as e.g. dependency failures are normal
# https://github.com/systemd/systemd/issues/12352
ExecStopPost=/usr/libexec/gnome-session-ctl --exec-stop-check
